# bootstrap/droid-config-pipa/sparse/etc/systemd/system/usb-moded.service.d/10-pipa.conf
[Unit]
Wants=
Requires=setup-configfs.service
After=local-fs.target dbus.socket sys-kernel-config.mount setup-configfs.service systemd-modules-load.service
After=multi-user.target

[Service]
# -r: enable developer_mode networking without android power_supply/usb
# (mainline has no /sys/class/power_supply/usb — usb-moded would never softconnect)
Environment=USB_MODED_ARGS=-r
# Ensure UDC unbound so idVendor writes succeed
ExecStartPre=/bin/sh -c 'echo "" > /sys/kernel/config/usb_gadget/g1/UDC 2>/dev/null || true'
ExecStartPre=/bin/sh -c 'test -e /sys/kernel/config/usb_gadget/g1/UDC || /usr/bin/droid/setup-configfs.sh'
TimeoutStartSec=30
