# droid-config-pipa/sparse/etc/systemd/system/start-user-session.service.d/10-pipa-timeout.conf
[Service]
# Stock uses TimeoutStartSec=0 (forever). Cap so lipstick failures do not stall boot.
TimeoutStartSec=120
