# bootstrap/droid-config-pipa/sparse/usr/lib/systemd/system/sensorfwd.service.d/ssc.conf
[Unit]
After=hexagonrpcd-sdsp.service mnt-vendor-persist.mount
Wants=hexagonrpcd-sdsp.service mnt-vendor-persist.mount
