# droid-config-pipa/sparse/usr/lib/systemd/system/pipa-rfkill-unblock.service
[Unit]
Description=Unblock WiFi/BT rfkill on pipa
After=systemd-modules-load.service
[Service]
Type=oneshot
ExecStart=/usr/sbin/rfkill unblock all
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
