# droid-config-pipa/sparse/usr/lib/systemd/system/pipa-bt-reprobe.service
[Unit]
Description=Reprobe QCA6390 Bluetooth after firmware is available
After=local-fs.target systemd-modules-load.service
Before=bluetooth.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for d in /sys/bus/serial/devices/*/uevent; do grep -q bluetooth "$(dirname $d)/of_node/compatible" 2>/dev/null || continue; echo add > "$d"; done; sleep 1; true'
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
