# bootstrap/droid-config-pipa/sparse/usr/lib/systemd/system/pipa-usb-firewall-allow.service
[Unit]
Description=Allow USB RNDIS (172.16.42.0/24) through ConnMan firewall
After=connman.service usb-rndis.service network-online.target
Wants=connman.service

[Service]
Type=oneshot
RemainAfterExit=yes
# ConnMan rebuilds iptables after start; re-open developer USB each time.
ExecStart=/bin/sh -c 'sleep 2; for i in usb0 rndis0; do iptables -C INPUT -i "$i" -j ACCEPT 2>/dev/null || iptables -I INPUT -i "$i" -j ACCEPT; iptables -C OUTPUT -o "$i" -j ACCEPT 2>/dev/null || iptables -I OUTPUT -o "$i" -j ACCEPT; done; iptables -C INPUT -s 172.16.42.0/24 -j ACCEPT 2>/dev/null || iptables -I INPUT -s 172.16.42.0/24 -j ACCEPT; iptables -C OUTPUT -d 172.16.42.0/24 -j ACCEPT 2>/dev/null || iptables -I OUTPUT -d 172.16.42.0/24 -j ACCEPT; true'

[Install]
WantedBy=multi-user.target
