# bootstrap/droid-config-pipa/sparse/usr/lib/systemd/system/usb-rndis.service
[Unit]
Description=USB RNDIS gadget for Sailfish pipa SSH
Conflicts=usb-moded.service
After=local-fs.target systemd-modules-load.service
Wants=systemd-modules-load.service
# Start with multi-user only — not sysinit (avoids early-boot stalls).

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=20
ExecStartPre=-/bin/systemctl stop usb-moded.service
ExecStart=/usr/bin/usb-rndis-gadget.sh
ExecStop=/bin/bash -c 'echo "" > /sys/kernel/config/usb_gadget/g1/UDC 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
