# bootstrap/droid-config-pipa/sparse/etc/profile.d/qt-camerabin.sh
# Sailfish Qt camerabin + libcamera (pipa)
export QT_GSTREAMER_CAMERABIN_SRC=wrappercamerabinsrc
export QT_GSTREAMER_CAMERABIN_VIDEOSRC=v4l2src
export QT_GSTREAMER_CAMERABIN_FLAGS=15
export QT_GSTREAMER_PLAYBIN_FLAGS=0
