# droid-config-pipa/sparse/usr/lib/systemd/system/connman.service.d/wlan-mac.conf
[Unit]
# WiFi MAC must be stable before ConnMan indexes wifi_<mac>_... services
After=wlan-module-load.service
Wants=wlan-module-load.service
