# droid-config-pipa/sparse/usr/lib/systemd/system/hexagonrpcd-sdsp.service.d/persist.conf
[Unit]
After=mnt-vendor-persist.mount
Requires=mnt-vendor-persist.mount
